# roles/GCBBase/files/servers/serverName1/etc/systemd/system/local-fs.target.wants/tmp.mount
# TWGCB-01-008-0004
[Mount]
Options=mode=1777,strictatime,noexec,nodev,nosuid

# TWGCB-01-008-0005
[Mount]
What=tmpfs
Where=/tmp
Type=tmpfs

